# supervisor.ini.j2
[program:{{ PROGRAM }}]
autostart=true
autorestart=true
stopsignal=QUIT
user=vagrant
command={{ COMMAND }}
environment=SETTINGS={{ ENVIRONMENT }}
