# provision_devenv/deed-api-stub.ini
[program:deed_api_stub]
autostart=true
autorestart=true
stopsignal=QUIT
user=vagrant
command=deed_api_stub -p 5004
environment=SETTINGS=deed_api_stub.config.Config